# commit-farming-service-template.service
[Unit]
Description="Commit farming script"

[Service]
ExecStart=/root/commit-farming/commit-farming.sh
